# infra/packer/ops/systemd/platform-update.timer
[Unit]
Description=Run platform update check every 5 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
