# infra/packer/ops/systemd/platform-update.service
[Unit]
Description=Check for new platform API image

[Service]
Type=oneshot
ExecStart=/usr/local/bin/platform-update.sh
